# files/etc/supervisor.d/dnsmasq.ini
[program:dnsmasq]
command=/usr/sbin/dnsmasq -k -C /etc/dnsmasq.conf
autostart=true
autorestart=true
priority=10