# router/etc/supervisord.d/named.ini
[program:named]
command=/usr/sbin/named -c /etc/named.conf -f -u named
autostart=true
autorestart=true
priority=50